# debian.default/puavo-ltsp-client.puavo-update-certificate.cron.daily
#!/bin/sh

set -eu

/usr/sbin/puavo-update-certificate
